# infrastructure/nexus-update.service
[Unit]
Description=Project Nexus Automatic Update Service
After=network-online.target postgresql.service neo4j.service
Wants=network-online.target

[Service]
Type=oneshot
User=nexus
Group=nexus
ExecStart=${pkgs.bash}/bin/bash -c 'export PATH="/opt/nexus/backend/node_modules/.bin:$PATH" && /opt/nexus/scripts/update-nexus.sh --cron'
StandardOutput=journal
StandardError=journal
